# slurm/s6-rc.d/munge/run
#!/bin/sh
cp /secrets/munge.key /etc/munge/munge.key || echo "Copy file has failed"
chown munge: /etc/munge/munge.key
s6-setuidgid munge:munge /usr/sbin/munged -F
