# slurm/s6-rc.d/munge/run
#!/bin/sh
set -ex

cp /secrets/munge/munge.key /etc/munge/munge.key
chown -R munge: /etc/munge /var/lib/munge /var/log/munge /var/run/munge
s6-setuidgid munge /usr/sbin/munged -F
