# slurm/s6-rc.d/munge/run
#!/bin/sh
cp /secrets/munge.key /etc/munge/munge.key || echo "Copy file has failed"
chown munge: /etc/munge/munge.key
/usr/sbin/munged -F
