# slurm/s6-rc.d/slurmrest/run
#!/bin/bash
set -ex

export SLURM_JWT=daemon
s6-applyuidgid -u "$(id api -u)" -g "$(id api -g)" /usr/sbin/slurmrestd 127.0.0.1:6820
