# slurm/s6-rc.d/oddjob/run
#!/bin/sh
set -ex

# -d no daemon, -n debug
/sbin/oddjobd -d -n
