# slurm/s6-rc.d/slurm/run
#!/bin/bash
set -ex
cp /secrets/slurm/jwt_hs256.key /etc/slurm/jwt_hs256.key
while ! id slurm
do
  sleep 1
done
chown slurm:slurm /etc/slurm/jwt_hs256.key
chmod 0600 /etc/slurm/jwt_hs256.key
chown root:root /etc/slurm
chmod 0755 /etc/slurm
chown -R slurm: /var/spool/slurmd /var/run/slurmd
/usr/sbin/slurmd -D
