# slurm/s6-rc.d/slurm/run
#!/bin/bash
set -x
while ! id slurm; do
  sleep 1
done
chown root:root /etc/slurm
chmod 0755 /etc/slurm
mkdir -p /var/spool/slurmd /var/run/slurmd /var/log/slurm
chown -R slurm: /var/spool/slurmd /var/run/slurmd
set -e
/usr/sbin/slurmd -D
