# slurm/s6-rc.d/sss/run
#!/bin/sh
set -ex

# -d no daemon, -n debug
/usr/sbin/sssd -d -n
