# slurm/s6-rc.d/sss/run
#!/bin/sh
set -ex

cp /secrets/sssd/sssd.conf /etc/sssd/sssd.conf
chmod 400 /etc/sssd/sssd.conf
chown root:root /etc/sssd/sssd.conf
/usr/sbin/sssd -i
